# config/cron.d/python_local_package_index_mirror
45 5 * * * /path/pep381client/pep381run -q /path/pypi/data
